# docs/deploy/status-collector.service
[Unit]
Description=Status Hub Collector
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/status-hub/status-collector.env
ExecStart=/usr/local/bin/status-collector -hub ${STATUS_HUB_URL} -token ${STATUS_DEVICE_TOKEN} -name ${STATUS_DEVICE_NAME}
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
